# ptools-daemon.service
[Unit]
Description=Manage ns-3-dce instances in background to provide better results

[Service]
Restart=on-failure
ExecStart=/usr/bin/python3.7 /root/playground/ptools/ptools_daemon.py
WorkingDirectory=/root/playground/ptools
RuntimeDirectoryMode=0777

[Install]
WantedBy=multi-user.target
